# systemd/Morfeas_system.service.d/Morfeas_system.conf
[Service]
#Add absolute path for the Morfeas_deamon configuration(.xml) file
Environment="Morfeas_deamon_config=Add_path_here"
